# roles/img-panda/templates/img-panda.conf
description "Img Panda Service"

start on runlevel [2345]
stop on runlevel [!2345]

exec {{ service_dir }}/img-panda/app.py >> {{ log_file }} 2>&1
